# .opencode/skills/kiro/com.kiro.webhook-receiver.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kiro.webhook-receiver</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/kokayi/Dev/claude-agent-workflow/.opencode/skills/kiro/webhook_receiver.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kokayi/Dev/claude-agent-workflow</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/kokayi/Library/Logs/kiro-webhook.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kokayi/Library/Logs/kiro-webhook.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
