# .coo/agent/launchd/com.yourcompany.coo.twitter-thread.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- TODO: Update label to match your company name -->
    <key>Label</key>
    <string>com.yourcompany.coo.twitter-thread</string>

    <key>ProgramArguments</key>
    <array>
        <!-- TODO: Update path to your project -->
        <string>/path/to/your/project/.coo/agent/launchd/run-scheduled-task.sh</string>
        <string>twitter-thread</string>
    </array>

    <!-- Runs at 4:30 PM on Tue (2) and Thu (4) -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>4</integer>
            <key>Hour</key>
            <integer>16</integer>
            <key>Minute</key>
            <integer>30</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/coo-agent/twitter-thread-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/coo-agent/twitter-thread-stderr.log</string>

    <!-- TODO: Update path to your project -->
    <key>WorkingDirectory</key>
    <string>/path/to/your/project</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
